# packaging/scripts/postinstall.sh
systemctl daemon-reload
systemctl enable slo-adapter.service
systemctl restart slo-adapter.service
