# packaging/scripts/preremove.sh
systemctl stop slo-adapter.service
systemctl disable slo-adapter.service
